# images/dind_prepare.service
[Unit]
Description=Prepare to start docker-in-docker

[Service]
Type=oneshot
ExecStart=/usr/sbin/dind_prepare

[Install]
RequiredBy=docker.service
